# hermes.service
[Unit]
Description=Hermes AI Agent
After=network.target
Documentation=https://github.com/NousResearch/hermes-agent
# Restart throttle (StartLimit* belong in [Unit] on systemd >= v230)
StartLimitIntervalSec=60s
StartLimitBurst=3

[Service]
Type=simple
User=hermes
WorkingDirectory=/opt/hermes
ExecStart=/opt/hermes/start.sh
EnvironmentFile=/etc/hermes.env

# systemd does NOT source /etc/profile.d — set the runtime PATH here so start.sh
# finds `hermes` (venv), `opencode` (npm-global), and node/python/uv/bun (mise shims).
Environment=PATH=/opt/hermes/.venv/bin:/opt/hermes/npm-global/bin:/opt/mise/shims:/usr/local/bin:/usr/bin:/bin
# Static runtime env from Dockerfile ENV block (lines 79-84)
Environment=HERMES_APP_DIR=/opt/hermes
Environment=HERMES_HOME=/opt/data
Environment=HERMES_WEBUI_REPO=/opt/hermes-webui
Environment=HERMES_WEBUI_TRUST_FORWARDED_HOST=1
Environment=PYTHONUNBUFFERED=1
Environment=MISE_DATA_DIR=/opt/mise
Environment=HF_HUB_ENABLE_HF_TRANSFER=1
Environment=PLAYWRIGHT_CHROMIUM_EXECUTABLE_PATH=/usr/bin/chromium

# Graceful shutdown: SIGTERM allows start.sh's trap to run final sync
KillSignal=SIGTERM
KillMode=mixed
TimeoutStopSec=60

# Automatic restart on failure (start.sh has internal supervisor)
Restart=on-failure
RestartSec=5

# Standard output/error to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes

[Install]
WantedBy=multi-user.target
